# dist/mastodon-sidekiq-default-@.service
[Unit]
Description=mastodon-sidekiq-default-%i
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/mnt/sdb/mastodon/live
Environment="RAILS_ENV=production"
Environment="DB_POOL=50"
Environment="MALLOC_ARENA_MAX=2"
ExecStart=/home/mastodon/.rbenv/shims/bundle exec sidekiq -c 50 -q default
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
