# dist/mastodon-sidekiq-pull-@.service
[Unit]
Description=mastodon-sidekiq-pull-%i
After=network.target

[Service]
Type=simple
User=mastodon
WorkingDirectory=/mnt/sdb/mastodon/live
Environment="RAILS_ENV=production"
Environment="DB_POOL=25"
Environment="MALLOC_ARENA_MAX=2"
ExecStart=/home/mastodon/.rbenv/shims/bundle exec sidekiq -c 25 -q pull
TimeoutSec=15
Restart=always

[Install]
WantedBy=multi-user.target
